# crontab
*/30 * * * * python /script .py >> /logfile
